# spider/supervisord.conf
[program:grab_task]
command = python /home/bae/app/grab_task.py
numprocess = 1
autorestart = false
stdout_logfile = /home/bae/log/main.log
stderr_logfile = /home/bae/log/mainerror.log

